# default.service
[Unit]
Description=A simple startup script

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/default.sh

[Install]
WantedBy=multi-user.target
